# infra/botblock-firewall/botblock-watcher.service
[Unit]
Description=BotBlock Watcher - Instant firewall blocking for detected bots
Documentation=https://github.com/youruser/botblock-firewall
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/botblock-watcher
Restart=always
RestartSec=5
StandardOutput=append:/var/log/botblock-watcher.log
StandardError=append:/var/log/botblock-watcher.log

# Security hardening
NoNewPrivileges=no
ProtectSystem=strict
ReadWritePaths=/tmp

[Install]
WantedBy=multi-user.target
